# debian/bms_ota.service
[Unit]
Description=BMS Firmware OTA Update Task
After=network.target
# 确保升级时不会被守护进程干扰（脚本内部已处理，这里是双重保险）
Conflicts=bms.service

[Service]
Type=oneshot
# 默认寻找 /opt/bms/firmware.bin 文件进行升级
ExecStart=/usr/bin/python3 /usr/local/bin/bms_ota.py /etc/default/bms_daemon /opt/bms/firmware.bin
RemainAfterExit=no
User=root

[Install]
WantedBy=multi-user.target
