# debian/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable bms.service
systemctl restart bms.service || true
