# debian/prerm
#!/bin/bash
systemctl stop bms.service || true
systemctl disable bms.service || true
